# systemd_templates/pulp-cleaner.timer
[Unit]
Requires=pulp-cleaner.service

[Timer]
Unit=pulp-cleaner.service
OnCalendar=*-*-01

[Install]
WantedBy=timers.target
